# ib-gateway.service
[Unit]
Description=IB Gateway via IBC (paper trading)
After=network.target xvfb.service
Requires=xvfb.service

[Service]
Type=simple
User=ec2-user
Environment=DISPLAY=:1
ExecStart=/home/ec2-user/ibc/gatewaystart.sh -inline \
    --gateway \
    --mode=paper \
    --ibc_path=/home/ec2-user/ibc \
    --ibc_ini=/home/ec2-user/ibc/config.ini \
    --java_path=/usr/lib/jvm/java-11-amazon-corretto \
    --tws_path=/home/ec2-user/Jts
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
