# mnq-gateway-morning-restart.timer
[Unit]
Description=Morning IB Gateway restart — Mon-Fri 07:30 ET, 2h before market open

[Timer]
OnCalendar=Mon..Fri *-*-* 07:30:00 America/New_York
Persistent=true

[Install]
WantedBy=timers.target
