# mnq-alerts.timer
[Unit]
Description=Start MNQ Alert System at market open (9:30 AM ET, weekdays)

[Timer]
OnCalendar=Mon..Fri *-*-* 09:30:00
Persistent=true

[Install]
WantedBy=timers.target
