# mnq-gateway-morning-restart.timer
[Unit]
Description=Daily IB Gateway restart — 07:30 ET, clears overnight login failures

[Timer]
OnCalendar=*-*-* 07:30:00 America/New_York
Persistent=true

[Install]
WantedBy=timers.target
