# mnq-2fa-reminder.timer
[Unit]
Description=2FA reminder — Sun-Thu 8:30 PM PT (before each trading day)

[Timer]
OnCalendar=Sun *-*-* 20:30:00 America/Los_Angeles
OnCalendar=Mon *-*-* 20:30:00 America/Los_Angeles
OnCalendar=Tue *-*-* 20:30:00 America/Los_Angeles
OnCalendar=Wed *-*-* 20:30:00 America/Los_Angeles
OnCalendar=Thu *-*-* 20:30:00 America/Los_Angeles
Persistent=true

[Install]
WantedBy=timers.target
