# mnq-backup.timer
[Unit]
Description=Back up alerts_log.db to S3 after market close (4:15 PM ET, weekdays)

[Timer]
OnCalendar=Mon..Fri *-*-* 16:15:00
Persistent=true

[Install]
WantedBy=timers.target
