# mnq-backup.service
[Unit]
Description=Back up alerts_log.db and parquet files to S3

[Service]
Type=oneshot
User=ec2-user
EnvironmentFile=/home/ec2-user/future-of-trading-futures/mnq_alerts/.env
ExecStart=/home/ec2-user/future-of-trading-futures/mnq_alerts/scripts/backup.sh
