# mnq-2fa-reminder.timer
[Unit]
Description=Weekly 2FA reminder — Sunday 6 PM Pacific (9 PM ET)

[Timer]
OnCalendar=Sun *-*-* 21:00:00 America/New_York
Persistent=true

[Install]
WantedBy=timers.target
